# templates/default/filebeat/filebeat.xml.erb
<?xml version='1.0'?>

<!--
This file was generated by Chef for <%= node['fqdn'] %> with care.
Do not modify this file by hand!
-->

<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
Service manifest for filebeat
-->

<service_bundle type='manifest' name='filebeat:filebeat'>

<service
  name='<%= @name %>'
  type='service'
  version='1'>
  <create_default_instance enabled='false' />
  <single_instance />

  <dependency name='fs'
    grouping='require_all'
    restart_on='none'
    type='service'>
    <service_fmri value='svc:/system/filesystem/local' />
  </dependency>

  <dependency name='net'
    grouping='require_all'
    restart_on='none'
    type='service'>
    <service_fmri value='svc:/network/loopback' />
  </dependency>

  <exec_method
    type='method'
    name='start'
    exec='sh <%= @script %> start'
    timeout_seconds='120'>
  </exec_method>

  <exec_method
    type='method'
    name='stop'
    exec='sh <%= @script %> stop'
    timeout_seconds='60'>
  </exec_method>

</service>

</service_bundle>

